# utils/update_pbs_cache.timer
[Unit]
Description=cache pbs data every minute

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1m

[Install]
WantedBy=multi-user.target